# KEEP/etc/inittab
# Note that busybox uses a non-standard inittab.

::sysinit:/etc/rc.boot
::respawn:/bin/runsvdir-start
::shutdown:/etc/rc.shutdown
